# roles/custom-configs/files/rc.local
$ANSIBLE_VAULT;1.1;AES256
38343835323561373531363530656539643263363833623763333764663761343336653963613037
6661633335303638373266663239633539316264363661370a646464616364333666633932616362
31633735633966306238633135363534376439306630303161356230343062613634366632636362
3136353930613336360a653263333537383638623865656234343966386634323938383965343362
64636237313837373733333938393438396531313861396436653235633761636639373566336664
34363863613865323438316433333062636433303664323266613664373638643130663431643164
65336266613765653863653665626636326430306366353539376333646663336438346564313337
36656235343934643434313061336539343861633265393939313166376262343137336263666637
62373433386438333936616535383935313665663634636465643262373365376531663534306637
39626136373131366533376364653463323062636537666466326339636466333261326666353535
64363235306538313933653864383736626465393133353966396430353432363163666364356630
61623739303164353264343139303561393933326434346134623561633730353337363833393331
66323630636134633361333135373938663331313034643739303164643635643538623663313034
66616538346536666566383536363731613165633865376665613864346230666361363536316435
33643431663463613961656561313231366364396565343235346436363734336332663339616133
35373036366366333466386430346564396631343764363137383731636338373635653261306539
37393832326466303137353239393939636565616364333663653761636263323962313631336636
31363534313466626437336439643239633430396335373234656531326232393432623964346535
63666461663262303435343263653538333330383638323331613237366437633136303064373163
65623838626430393838336530303162313061613638613637393533663362306238636134303138
35623461653935313431663862356430383634653962663037633238316265616665613865366162
32336366393039623833383365313237373332373639376666363932626337373836663266373765
31323632666138656537333663303864653261336134623134316363313862666234346438613239
64613662643530663961613434373532353335393038343432613933346662613330633432346132
3730
